# scripts/policydraft.conf
[program:policydraft_backend]
command=python3 {{APP_DIR}}/backend/server.py
directory={{APP_DIR}}/backend
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/policydraft.log
stderr_logfile=/var/log/supervisor/policydraft.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
